# machines/relaycreator/deploy.timer
[Unit]
Description=Deploy relaycreator

[Timer]
# run checks once a minute
OnCalendar=*-*-* *:*:05

[Install]
WantedBy=timers.target