# dockerfiles/strfrycheck.timer
[Unit]
Description=Run strfrycheck

[Timer]
# run checks once a minute
OnCalendar=*-*-* *:*:05

[Install]
WantedBy=timers.target
